# ddns.service
[Unit]
Description=Checks for dynamic IP changes

[Service]
User=nobody
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/usr/bin/ddns
